# systemd/gearmand.service
#customization gearmand, requires mariadb - DGI
.include /usr/lib/systemd/system/gearmand.service
[Unit]
After=mariadb.service