# systemd/gearman-job-server.service
.include /lib/systemd/system/gearman-job-server.service

[Unit]
After=mysql.service
